# etc/runit/sv/scron/run
#!/bin/sh
exec /usr/local/bin/crond -n 2>&1
